# etc/systemd/system/medipick-mapping-app-bridge.service
[Unit]
Description=MediPick Mapping App Bridge
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
User=jszn
Group=jszn
WorkingDirectory=/home/jszn/robot_ros2_ws
Environment=HOME=/home/jszn
Environment=ROS_DISTRO=humble
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && source /home/jszn/robot_ros2_ws/install/setup.bash && exec ros2 launch robot_hardware mapping_app_bridge.launch.py device_name:=MediPick-Desk app_max_linear_x:=3.0 app_max_linear_y:=3.0 app_max_angular_z:=3.0'
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
